# blockstore/service/blockstore.tserver/run
#!/bin/bash

exec /usr/local/bin/run-service haobtc 'bin/blockstored -h 127.0.0.1'

